# files/docker/files/supervisord.conf
[supervisord]
nodaemon=true
user=minecraft
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:minecraft]
command=/start
